# example.crontab
# m h  dom mon dow   command
0 * * * * python digits.py
